# deploy/botshop-bot.service
[Unit]
Description=botshop telegram bot
After=network.target

[Service]
Type=simple
WorkingDirectory=/www/wwwroot/botshop
ExecStart=/usr/bin/python3.10 /www/wwwroot/botshop/haopubot.py
Restart=always
RestartSec=3
User=root
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
